# etc/cron.d/cephalix.jobs
15 6,14,22 * * * * root /usr/sbin/cephalix_get_status_of_institutes.sh &> /dev/null
